# config/root/etc/s6-overlay/s6-rc.d/bgutil-ytdlp-pot-provider/run
#!/command/with-contenv bash

set -e
cd /app/bgutil-ytdlp-pot-provider

exec s6-setuidgid app \
    ./node server/build/main.js -p 4406
