# config/root/etc/services.d/redis/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid redis \
    /usr/bin/redis-server /etc/redis/redis.conf
