# config/root/etc/s6-overlay/s6-rc.d/huey-filesystem/run
#!/command/with-contenv bash

exec nice -n "${TUBESYNC_NICE:-1}" s6-setuidgid app \
    /usr/bin/python3 /app/manage.py djangohuey \
    --queue filesystem
