# config/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bash

cd /

exec /usr/sbin/nginx
