# config/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

cd /

/usr/sbin/nginx
