# config/root/etc/s6-overlay/s6-rc.d/tubesync-worker/run
#!/command/with-contenv bash

exec s6-setuidgid app \
    /usr/bin/python3 /app/manage.py process_tasks
